# packaging/vibegram.service
[Unit]
Description=vibegram — Telegram supervision layer for coding agents
After=network.target

[Service]
Type=simple
User=vibegram
StateDirectory=vibegram
WorkingDirectory=/var/lib/vibegram
Environment=VIBEGRAM_WORK_ROOT=/var/lib/vibegram
Environment=VIBEGRAM_STATE_DIR=/var/lib/vibegram/state
EnvironmentFile=/etc/vibegram/env
ExecStart=/usr/local/bin/vibegram
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
